# etc/cron.d/mediaing-updatedb
MAILTO="acis-root@acis.ufl.edu"
LANG="en_US.UTF-8"
SHELL=/bin/bash
PATH=/opt/idb-backend/.env/bin:/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m  h  dom mon dow user  command
20   3  *   *   *   root  idigbio-ingestion --logfile 'mediaing-updatedb.log' -v mediaing updatedb --daily &> /dev/null
