# etc/systemd/system/idb-collect-stats.timer
[Unit]
Description="idb collect-stats hourly repeat"

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
