# etc/systemd/system/idigbio-ingestion-update-publisher-recordset.timer
[Unit]
Description="idigbio-ingestion update-publisher-recordset hourly repeat"

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
